# momentum_ml/deploy/momentum-api-watchdog.timer
[Unit]
Description=Kör Momentum API-vakthund periodiskt

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5s

[Install]
WantedBy=timers.target
