# momentum_ml/deploy/momentum-train.service
[Unit]
Description=Momentum ML – nattlig pipeline-körning (data, träning, backtest)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hannesb
WorkingDirectory=/opt/momentum/momentum_ml
Environment=PATH=/opt/momentum/venv/bin
# Lämna en kärna åt API:t/OS – justera vid behov för din Pi-modell.
Environment=MOMENTUM_TRAINING_THREADS=3
# Obuffrad stdout/stderr så journalctl visar utskrifter i rätt ordning/tid -
# annars kan parent- och subprocess-utskrifter framstå som hopklumpade.
Environment=PYTHONUNBUFFERED=1
# Helt universum (alla market-cap-nivåer inkl. Small/Micro/Nano). Momentum-
# premien sitter starkast i mindre, mindre analyserade bolag, så dessa hör hemma
# i universumet. Likviditetsfiltret (config.UNIVERSE_MIN_AVG_TURNOVER) trimmar
# bort de som faktiskt inte går att handla.
#
# VARNING – RAM: hela Sverige-listan (~656 tickers) drar betydligt mer minne än
# Large/Mid (~108). På en 2GB-Pi med MemoryMax=1400M nedan finns risk för
# OOM-kill. Om träningen dör (se journal / results/health.log):
#   1) höj --min-turnover (färre, mer likvida bolag), och/eller
#   2) begränsa nivåer igen, t.ex. --market-cap "Large Cap" "Mid Cap" "Small Cap"
#      (Small Cap ger merparten av momentum-vinsten; Micro/Nano är mest illikvida
#       + värst survivorship-bias), och/eller
#   3) höj MemoryMax/MemorySwapMax om Pi:n har mer RAM, eller byt till större Pi.
# Behåll --skip-lstm: LSTM ovanpå hela universumet spränger minnet säkert.
ExecStart=/opt/momentum/venv/bin/python main.py --universe sweden --skip-lstm --ta-filter score

# Låg prioritet så API:t/dashboarden inte blir lidande under träningen.
Nice=15
IOSchedulingClass=idle
CPUWeight=20
# Hård minnesgräns. MÅSTE vara LÄGRE än fysiskt RAM (2GB-Pi har ~1.8Gi
# tillgängligt) – annars dödar kärnan hela systemet (hård omstart) istället
# för att bara döda träningsprocessen när minnet tar slut. 1.4G lämnar
# headroom åt OS + API. Sänk vid OOM-kills, men minska hellre arbetsmängden
# (--market-cap eller --skip-lstm) så taket inte ens nås.
MemoryMax=1400M
MemorySwapMax=512M

# Tidsgräns som extra säkerhetsnät mot en hängande körning (justera vid behov).
TimeoutStartSec=4h
