# momentum_ml/deploy/momentum-api-watchdog.service
[Unit]
Description=Momentum API health watchdog (restart if /api/health is unhealthy)
After=momentum-api.service

[Service]
Type=oneshot
# Körs som root (ingen User=) så att systemctl restart får behörighet.
ExecStart=/opt/momentum/momentum_ml/deploy/api_watchdog.sh
